# compiled/splunk/rules.conf
===============================================================================
   DETECTION RULES - COMPILED FOR SPLUNK
   Generated: 2026-03-21
   Source: Sigma Rules
===============================================================================

-------------------------------------------------------------------------------
Rule 1: Suspicious PowerShell Encoded Command Execution
Sigma ID: a1b2c3d4-e5f6-7890-abcd-ef1234567890
ATT&CK: T1059.001, T1027.010
Level: Critical
-------------------------------------------------------------------------------

index=windows sourcetype=WinEventLog:Sysmon EventCode=1
  (ParentImage="*\\cmd.exe" OR ParentImage="*\\wscript.exe"
   OR ParentImage="*\\cscript.exe" OR ParentImage="*\\mshta.exe"
   OR ParentImage="*\\wmiprvse.exe" OR ParentImage="*\\rundll32.exe")
  (Image="*\\powershell.exe" OR Image="*\\pwsh.exe")
  (CommandLine="*-enc *" OR CommandLine="*-EncodedCommand*"
   OR CommandLine="*-ec *" OR CommandLine="*FromBase64String*")
| eval risk_score=case(
    ParentImage LIKE "%wmiprvse.exe%", 90,
    ParentImage LIKE "%mshta.exe%", 85,
    ParentImage LIKE "%rundll32.exe%", 80,
    1=1, 70
  )
| where NOT (match(CommandLine, "(?i)SCCM") OR match(CommandLine, "(?i)ConfigMgr"))
| table _time Computer User ParentImage Image CommandLine risk_score ProcessId ParentProcessId
| sort - risk_score _time

-------------------------------------------------------------------------------
Rule 2: LSASS Memory Dumping Attempt
Sigma ID: b3c4d5e6-f7g8-9012-cdef-1234567890ab
ATT&CK: T1003.001, T1003.003
Level: Critical
-------------------------------------------------------------------------------

index=windows sourcetype=WinEventLog:Sysmon EventCode=10
  TargetImage="*\\lsass.exe"
  GrantedAccess IN ("0x1010", "0x1038", "0x1fffff", "0x1410")
| where NOT (
    SourceImage="C:\\Windows\\System32\\wininit.exe"
    OR SourceImage="C:\\Windows\\System32\\csrss.exe"
    OR SourceImage="C:\\Windows\\System32\\lsm.exe"
    OR match(SourceImage, "C:\\ProgramData\\Microsoft\\Windows Defender\\Platform\\.*\\MsMpEng.exe")
    OR match(SourceImage, "C:\\Program Files\\WindowsApps\\.*WindowsDefender\\.*\\MsMpEng.exe")
  )
| eval risk_score=90
| table _time Computer User SourceImage TargetImage GrantedAccess risk_score AccessMask ProcessId SourceProcessId
| sort - risk_score _time

-------------------------------------------------------------------------------
Rule 3: Process Hollowing / Injection Detection
Sigma ID: d4e5f6g7-h8i9-1234-efab-345678901cde
ATT&CK: T1055.012, T1055
Level: High
-------------------------------------------------------------------------------

index=windows sourcetype=WinEventLog:Sysmon EventCode=8
  NOT (SourceImage="C:\\Windows\\System32\\csrss.exe"
       OR SourceImage="C:\\Windows\\System32\\wininit.exe")
| where NOT (
    TargetImage="C:\\Windows\\System32\\csrss.exe"
  )
| eval risk_score=75
| table _time Computer User SourceImage TargetImage NewThreadId risk_score StartAddress SourceProcessId TargetProcessId
| sort - risk_score _time

-------------------------------------------------------------------------------
Rule 4: Registry Persistence via Run Keys
Sigma ID: f5g6h7i8-j9k0-3456-bcde-567890123ef0
ATT&CK: T1547.001
Level: Medium
-------------------------------------------------------------------------------

index=windows sourcetype=WinEventLog:Sysmon (EventCode=12 OR EventCode=13 OR EventCode=14)
  (TargetObject="*\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run*"
   OR TargetObject="*\\SOFTWARE\\Microsoft\\Windows NT\\CurrentVersion\\Winlogon*")
  (Details="*.exe" OR Details="*.bat" OR Details="*.cmd"
   OR Details="*.ps1" OR Details="*.vbs" OR Details="*.js")
| where NOT (
    match(Details, "^C:\\\\Program Files\\\\")
    OR match(Details, "^C:\\\\Program Files \\(x86\\\\)")
    OR match(Details, "^C:\\\\Windows\\\\System32\\\\")
    OR match(Details, "^C:\\\\Windows\\\\SysWOW64\\\\")
  )
| where NOT (
    match(TargetObject, "OneDriveSetup.exe")
    OR match(TargetObject, "MicrosoftEdgeAutoLaunch")
    OR match(Details, "OneDrive\\.exe")
    OR match(Details, "Razer Synapse")
  )
| eval risk_score=60
| table _time Computer User TargetObject Details risk_score ProcessId Image
| sort - risk_score _time

-------------------------------------------------------------------------------
Rule 5: Suspicious Script Network Connection
Sigma ID: g7h8i9j0-k1l2-5678-defa-789012345012
ATT&CK: T1071, T1059.001, T1059.005
Level: High
-------------------------------------------------------------------------------

index=windows sourcetype=WinEventLog:Sysmon EventCode=3
  (Image="*\\powershell.exe" OR Image="*\\pwsh.exe"
   OR Image="*\\wscript.exe" OR Image="*\\cscript.exe"
   OR Image="*\\mshta.exe" OR Image="*\\rundll32.exe"
   OR Image="*\\regsvr32.exe")
  NOT (DestinationIp="127.0.0.1" OR DestinationIp="0.0.0.0" OR DestinationIp="::1")
  NOT (
    match(DestinationIp, "^192\\.168\\.")
    OR match(DestinationIp, "^10\\.")
    OR match(DestinationIp, "^172\\.(1[6-9]|2[0-9]|3[0-1])\\.")
  )
| where NOT (
    match(DestinationHostname, "\\.microsoft\\.com$")
    OR match(DestinationHostname, "\\.windowsupdate\\.com$")
    OR match(DestinationHostname, "\\.windows\\.net$")
  )
| eval risk_score=75
| table _time Computer User Image DestinationIp DestinationPort DestinationHostname risk_score Protocol ProcessId
| sort - risk_score _time

===============================================================================
END OF SPLUNK RULES
===============================================================================
